# docker/etc/supervisor/conf.d/celery-beat.conf
[program:celery-beat]
command = celery -A project beat --loglevel=INFO
directory = /var/app/
stdout_logfile = /var/log/app-celery-beat.log
stderr_logfile =/var/log/app-celery-beat.log
user = root
autorestart  = True
